# CloverPackage/CloverV2/etc/rc.local
#!/bin/bash
cat boot.log >boot.lo1
echo " " >> boot.lo1
date "+DATE: %Y-%m-%d TIME: %H:%M:%S%n" >> boot.lo1
ioreg -lw0 -pIODeviceTree | grep boot-log > ./log.txt
logdump=$(cat ./log.txt)
modified1=${logdump#*'boot-log'}
modified2=${modified1#*'<'}
modified3=${modified2%%'>'*}
echo $modified3 > ./boot-log.txt
xxd -r -p ./boot-log.txt >> boot.lo1
rm log.txt
rm boot-log.txt
tail -n 3000 boot.lo1 >boot.log
rm boot.lo1
mkdir /Volumes/EFI
mount_msdos /dev/disk0s1 /Volumes/EFI