# CloverPackage/CloverV2/rcScripts/Library/LaunchDaemons/com.projectosx.clover.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.projectosx.clover.daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/Clover/CloverDaemon</string>
		<string>startup</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Library/Logs/CloverEFI/clover.daemon.log</string>
	<key>StandardOutPath</key>
	<string>/Library/Logs/CloverEFI/clover.daemon.log</string>
	<key>WorkingDirectory</key>
	<string>/etc</string>
</dict>
</plist>
